# runit/fs.dbwriter/run
#!/bin/sh

exec 2>&1

USER='www'
DROPROOT="chpst -u ${USER}"

cd /www/filestorage

exec ${DROPROOT} php -c etc/php.ini dbwriter.php
